# crawler/com.xhs.dashboard.crawler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xhs.dashboard.crawler</string>

    <!-- setup_autostart.sh renders these placeholders before installing. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>MEDIACRAWLER_DIR</key>
        <string>__MEDIACRAWLER_DIR__</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>__CRAWLER_DIR__/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__MEDIACRAWLER_DIR__</string>

    <!-- 登录后自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃后自动重启 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 最小重启间隔30秒，防止崩溃循环 -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__CRAWLER_DIR__/logs/server-out.log</string>

    <key>StandardErrorPath</key>
    <string>__CRAWLER_DIR__/logs/server-err.log</string>
</dict>
</plist>
